# wherecron_data/crontabs/root
# crontab running on monday and thursday at 7:30
50 4 * * 1  /scripts/run_argos.sh >> /var/log/cron.log 2>&1
50 4 * * 4  /scripts/run_argos.sh >> /var/log/cron.log 2>&1

